# deploy/face-attendance.service
[Unit]
Description=Face recognition attendance (Gunicorn + Flask)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DEPLOY_USER__
Group=__DEPLOY_USER__
WorkingDirectory=__INSTALL_ROOT__
EnvironmentFile=-/etc/face-attendance.env

# Port from FLASK_PORT in EnvironmentFile (default 8002 in gunicorn_start.sh).
ExecStart=__INSTALL_ROOT__/deploy/gunicorn_start.sh

Restart=on-failure
RestartSec=5

# Hardening (optional; comment out if you hit permission issues reading camera-less paths)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
